# .crontab
#!/bin/bash

*/5 * * * * curl -X GET http://chronos-edt.herokuapp.com
